# infra/ucore/containers/minio.container
# Podman Quadlet container definition for MinIO
# Place in /etc/containers/systemd/ or include in Ignition config

[Unit]
Description=MinIO S3-compatible object storage
After=network-online.target zfs-import-tank.service
Wants=network-online.target

[Container]
Image=quay.io/minio/minio:latest
ContainerName=minio
PublishPort=9000:9000
PublishPort=9001:9001

# Mount ZFS paths
Volume=/var/tank/services/minio/cfg:/root/.minio:Z
Volume=/var/tank/services/minio/data:/data:Z

# MinIO configuration
Environment=MINIO_ROOT_USER=admin
Environment=MINIO_ROOT_PASSWORD=changeme
Environment=MINIO_REGION=us-kyz-0

# Command
Exec=server /data --console-address ":9001"

# Auto-restart
Restart=on-failure
RestartSec=10s

[Service]
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
